# CITIUS/Pruebas/PruebasJAUS/JAUS_Subs2_Node2/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/atica/catkin_ws/src/CITIUS/Pruebas/PruebasJAUS/JAUS_Subs2_Node2/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
